# server/root/etc/services.d/pagermon/run
#!/usr/bin/with-contenv bash

if [ -n "$SKIP_APP" ]; then
	exit 0
fi

cd /app

exec s6-setuidgid $CONTAINER_USER node app.js
